# tema5/listado_curso.timer
[Unit]
Description="Para hacer listado de curso cada 20min"

[Timer]
OnCalendar=*:0/20

[Install]
WantedBy=timers.target

