# tema5/revision.service
[Unit]
Description="Script para revisar el registro syslog"

[Service]
Type=oneshot
ExecStart=/home/jefe/scripts/revisar.sh
